# config/puma.service
[Unit]
Description=Puma HTTP Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/robo-catcher
ExecStart=/bin/bash -lc 'bundle exec puma -C /home/pi/robo-catcher/config/puma.rb PORT=80'
User=pi
Group=pi
Restart=always
Environment=RAILS_ENV=production

[Socket]
ListenStream=0.0.0.0:80
NoDelay=true
ReusePort=true
Backlog=1024
Restart=always

[Install]
WantedBy=sockets.target
