# config/puma.service
[Unit]
Description=Rails Server
Requires=network.target

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/robo-catcher
ExecStart=/home/pi/robo-catcher/bin/bundle exec rails server -e production
TimeoutSec=30
RestartSec=15s
Restart=always

[Install]
WantedBy=multi-user.target
